# infra/systemd/bilibili-views-crawler-cron.timer
[Unit]
Description=B站投稿数据爬虫定时器 (Cron脚本版)

[Timer]
# 每小时执行一次
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
