# infra/systemd/moments-crawler.service
[Unit]
Description=满の动态爬虫服务（微博+B站动态）
Documentation=https://github.com/xxm-fans/xxm_fans_home
After=network.target

[Service]
Type=oneshot
User=yifeianyi
Group=yifeianyi

WorkingDirectory=/home/yifeianyi/Desktop/xxm_fans_home
ExecStart=/home/yifeianyi/Desktop/xxm_fans_home/scripts/moments_cron.sh

Environment="PYTHONUNBUFFERED=1"
Environment="DJANGO_SETTINGS_MODULE=xxm_fans_home.settings"

MemoryLimit=512M
CPUQuota=50%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=moments-crawler

Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
